# chat-relay.service
[Unit]
Description=Chat Relay Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/chat-relay
ExecStart=/usr/bin/node /home/ubuntu/chat-relay/dist/index.js
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/chat-relay/logs/systemd.log
StandardError=append:/home/ubuntu/chat-relay/logs/systemd-error.log
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target